# backend/laravel-scheduler.service
[Unit]
Description=Laravel Scheduler for MIS Asset Management System
Documentation=https://laravel.com/docs/scheduling
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/mis-system/backend
ExecStart=/usr/bin/php artisan schedule:work
Restart=always
RestartSec=3

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=append:/var/www/mis-system/backend/storage/logs/scheduler.log
StandardError=append:/var/www/mis-system/backend/storage/logs/scheduler-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
